# send-attachments.crontab
* * * * * root /opt/send-attachments.sh >> /var/log/send-attachments.log 2>&1

